# rygel/supervisor.conf
[program:rygel]
command=rygel -c /etc/rygel.conf
user={{ pillar.media_server.user }}
autostart=True
autorestart=True
redirect_stderr=True
stopasgroup=True
killasgroup=True

